# cgminer/files/root.cron
*/3 *  *   *   *     /usr/bin/cgminer-monitor